# supervisor/horizon-sample.conf

[program:laravel-horizon]
process_name=%(program_name)s
command=/usr/bin/php /var/www/html/artisan horizon
numprocs=1
autostart=true
autorestart=true
